# wifiphisher/rootfs/etc/services.d/wifiphisher/finish
#!/bin/sh
echo "[wifiphisher] service stopped"
