# wifiphisher/rootfs/etc/services.d/wifiphisher/run
#!/usr/bin/with-contenv bash
set -e

echo "[wifiphisher] starting..."

# Falls Interface nicht gesetzt ist
IFACE="${INTERFACE:-wlan0}"

# Warten bis Interface existiert (wichtig!)
while ! ip link show "$IFACE" >/dev/null 2>&1; do
  echo "[wifiphisher] waiting for interface $IFACE..."
  sleep 2
done

echo "[wifiphisher] using interface: $IFACE"

# Starte wifiphisher im Vordergrund (wichtig!)
exec wifiphisher --interface "$IFACE"
